# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

start-stop-daemon -b -S -x datareport.lua
start-stop-daemon -b -S -x commanddispatcher.lua
start-stop-daemon -b -S -x flowmeter
start-stop-daemon -b -S -x ttyswitch
sleep 1
start-stop-daemon -b -S -x zlldaemon
<<<<<<< HEAD
sleep 1
=======
sleep 2
>>>>>>> ba1ea7de21aceb28cd82671090139b1f33cbc006
start-stop-daemon -b -S -x netzigbee
start-stop-daemon -b -S -x webdevice.lua
start-stop-daemon -b -S -x wifibridge
start-stop-daemon -b -S -x wifibridge_1

exit 0
