# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

start-stop-daemon -b -S -x datareport.lua
start-stop-daemon -b -S -x commanddispatcher.lua
start-stop-daemon -b -S -x flowmeter
start-stop-daemon -b -S -x ttyswitch
sleep 1
start-stop-daemon -b -S -x /mnt/zll -- /dev/ttyS13
start-stop-daemon -b -S -x netzigbee
start-stop-daemon -b -S -x webdevice.lua


exit 0
